# services/com.aya.embedding-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aya.embedding-service</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Volumes/DATA/AYA/services/embedding_service.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/DATA/AYA/services</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/arthurdell/Library/Logs/AgentTurbo/embedding.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/arthurdell/Library/Logs/AgentTurbo/embedding_error.log</string>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>

