# scripts/com.aya.git.autosync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aya.git.autosync</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Volumes/DATA/AYA/scripts/auto_git_sync.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Volumes/DATA/AYA</string>
    
    <!-- Run every 15 minutes -->
    <key>StartInterval</key>
    <integer>900</integer>
    
    <!-- Also run at load -->
    <key>RunAtLoad</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/Volumes/DATA/AYA/logs/git_autosync.out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Volumes/DATA/AYA/logs/git_autosync.err.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Library/Apple/usr/bin</string>
        <key>HOME</key>
        <string>/Users/arthurdell</string>
        <key>USER</key>
        <string>arthurdell</string>
    </dict>
    
    <!-- Don't run if network is unavailable -->
    <key>StartOnMount</key>
    <false/>
    
    <!-- Low priority, don't interrupt work -->
    <key>Nice</key>
    <integer>10</integer>
    
    <key>ThrottleInterval</key>
    <integer>60</integer>
</dict>
</plist>

